# backend/moneyone-api-improved.conf
[program:moneyone-api]
command=/var/www/moneyone/moneyone/backend/venv/bin/gunicorn --bind 127.0.0.1:5000 --workers 4 --timeout 120 --log-level info --access-logfile /var/log/moneyone/api-access.log --error-logfile /var/log/moneyone/api-error.log app:app
directory=/var/www/moneyone/moneyone/backend
user=ubuntu
autostart=true
autorestart=true
stopasgroup=true
killasgroup=true
stderr_logfile=/var/log/moneyone/api-stderr.log
stdout_logfile=/var/log/moneyone/api-stdout.log
environment=PATH="/var/www/moneyone/moneyone/backend/venv/bin"

; Logging configuration
stderr_logfile_maxbytes=50MB
stderr_logfile_backups=10
stdout_logfile_maxbytes=50MB
stdout_logfile_backups=10
redirect_stderr=false

; Process management
priority=999
startsecs=10
startretries=3
